# worker/supervisord.conf
[supervisord]
nodaemon=true

[program:rqworker]
command=python -u run-worker.py
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true



